# src/DEBIAN/preinst
#!/bin/sh

deb-systemd-helper stop apache2
deb-systemd-helper stop mariadb
